# root/etc/services.d/webapp/run
#!/usr/bin/with-contenv bash

cd /app/synclounge || exit

exec \
	s6-setuidgid abc node webapp.js
